# .pipeline/projects/email_tool/systemd/email_tool_daemon.service
[Unit]
Description=Email Tool Daemon - Background email synchronization service
Documentation=https://github.com/example/email-tool
After=network.target
Wants=network.target

[Service]
Type=forking
ExecStart=/usr/bin/python3 -m email_tool.daemon --foreground
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/usr/bin/python3 -m email_tool.daemon --stop
PIDFile=/tmp/email_tool_daemon.pid
Restart=on-failure
RestartSec=30

# Environment
Environment=EMAIL_TOOL_CONFIG=/etc/email_tool/config.yaml
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/email_tool
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true
LockPersonality=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=256

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=email_tool_daemon

[Install]
WantedBy=multi-user.target
