# .pipeline/projects/email_tool/workspace/examples/systemd/email-tool.service
[Unit]
Description=Email Tool Daemon Service
Documentation=https://github.com/email-tool/email-tool
After=network.target
Wants=network.target

[Service]
Type=simple
User=emailuser
Group=emailuser
WorkingDirectory=/opt/email-tool
ExecStart=/opt/email-tool/venv/bin/email-tool daemon start
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=email-tool

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/email-tool/data
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Resource limits
MemoryMax=512M
CPUQuota=80%

# Environment variables
Environment=EMAIL_TOOL_CONFIG=/opt/email-tool/config.yaml
Environment=LOG_LEVEL=INFO

[Install]
WantedBy=multi-user.target
