# .pipeline/projects/email_tool/workspace/examples/systemd/email-tool.timer
[Unit]
Description=Email Tool Scheduled Processing Timer
Documentation=https://github.com/email-tool/email-tool
Requires=email-tool.service
After=email-tool.service

[Timer]
# Run every 15 minutes
OnBootSec=5min
OnUnitActiveSec=15min

# Randomize start time to avoid thundering herd
RandomizedDelaySec=30

# Persist timer across reboots
Persistent=true

# Unit to activate
Unit=email-tool.service

[Install]
WantedBy=timers.target
